# systemd/offlinemsmtp.service
[Unit]
Description=offlinemsmtp

[Service]
ExecStart=/usr/bin/offlinemsmtp --daemon

[Install]
WantedBy=default.target
